# kde/files/config/autostart/disable-dpms.desktop
[Desktop Entry]
Type=Application
Name=Disable DPMS
Exec=xset -dpms
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
